# lib/goico/installer/templates/launchd.worker.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string><%= manifest["app_name"] %>.worker</string>
    <key>ProgramArguments</key>
    <array>
      <string><%= app_path %>/bin/rails</string>
      <string>jobs:work</string>
    </array>
    <key>WorkingDirectory</key>
    <string><%= app_path %></string>
    <key>UserName</key>
    <string><%= app_user %></string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
